# docker/supervisor/laravel-scheduler.conf
[supervisord]
nodaemon=true

[program:laravel-scheduler]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/artisan schedule:work
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www
redirect_stderr=true
stdout_logfile=/var/log/supervisor/scheduler.log
